# supervisor/conf/controller.ini
[program:controller]
command=/bin/su kafka /bin/bash -c "/opt/kafka/script/controller.sh"
process_name=%(program_name)s
pidfile=/var/run/controller.pid
numprocs=1
autorestart=false
autostart=true
startsecs=10
startretries=1
stopsignal=TERM
stopwaitsecs=10